# go/sentinel-agent.service
[Unit]
Description=Sentinel Monitoring Agent
Documentation=https://github.com/your-org/sentinel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

# Agent binary and configuration
ExecStart=/usr/local/bin/sentinel-agent -config /etc/sentinel-agent/config.yaml

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel-agent

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/sentinel-agent

# Resource limits
MemoryMax=128M
CPUQuota=10%

[Install]
WantedBy=multi-user.target
